# system_files/etc/profile.d/zmotd.sh
#!/usr/bin/env bash

/usr/libexec/zmotd
